# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:api-gateway]
command=npx tsx /app/services/api-gateway/src/index.ts
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3000,NODE_ENV=production,ORDER_SERVICE_URL="http://localhost:3001",DOCUMENT_SERVICE_URL="http://localhost:3002",NLP_SERVICE_URL="http://localhost:3003",INTEGRATION_SERVICE_URL="http://localhost:3004",CORS_ORIGIN="*"

[program:order-service]
command=npx tsx /app/services/order-service/src/index.ts
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3001,NODE_ENV=production

[program:document-service]
command=npx tsx /app/services/document-service/src/index.ts
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3002,NODE_ENV=production

[program:nlp-service]
command=node -e "const http=require('http');const s=http.createServer((q,r)=>{r.setHeader('Content-Type','application/json');if(q.url==='/health'){r.end(JSON.stringify({status:'healthy',service:'nlp-service',version:'1.0.0',note:'stub - Python NLP not deployed'}))}else{r.statusCode=404;r.end(JSON.stringify({error:'NLP service stub - full service not deployed'}))}});s.listen(3003,()=>console.log('NLP stub on 3003'))"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:integration-service]
command=npx tsx /app/services/integration-service/src/index.ts
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3004,NODE_ENV=production
